# deploy/audio-transcript-bot.service
[Unit]
Description=Audio Transcript Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/USER/audio-transcript-bot
ExecStart=/home/USER/audio-transcript-bot/.venv/bin/python bot.py
EnvironmentFile=/home/USER/audio-transcript-bot/.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
